# files/etc/uci-defaults/luci-setup-wizard
#!/bin/sh

. /etc/functions.sh

#setup ucitrack
config_load ucitrack
uci_add ucitrack nodogsplash nodogsplash
uci_set ucitrack nodogsplash init ucidog
uci_commit ucitrack

exit 0
