# files/etc/uci-defaults/ucitrack-updates
#!/bin/sh

. /etc/functions.sh

config_load ucitrack
uci_add ucitrack nodogsplash nodogsplash
uci_set ucitrack nodogsplash init ucidog
uci_commit ucitrack
/etc/init.d/ucidog enable
