# files/uci-defaults/commotion-dashboard-helper
#!/bin/sh
echo '* * * * * /usr/bin/commotion-bigboard-send' >> /etc/crontabs/root
